# ganglia_ramdisk/gmetad_checkpoint_ramdisk.cron
12 * * * * root /usr/local/sbin/gmetad_checkpoint_ramdisk.sh
